# com.vibecode.claude-bot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.vibecode.claude-bot</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/brianle/Documents/Code/Coding Tool/claude-code-telegram-boss/start_bot.sh</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/brianle/Documents/Code/Coding Tool/claude-code-telegram-boss/bot.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/brianle/Documents/Code/Coding Tool/claude-code-telegram-boss/bot.err</string>
    <key>WorkingDirectory</key>
    <string>/Users/brianle/Documents/Code/Coding Tool/claude-code-telegram-boss</string>
</dict>
</plist>
